# config/systemd/chrome.service
[Unit]
Description=Chrome Browser Service
After=x11vnc.service
Requires=x11vnc.service

[Service]
Type=simple
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/wilson_austin10/.Xauthority
Environment=HOME=/home/wilson_austin10
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/1000/bus
Environment=XDG_RUNTIME_DIR=/run/user/1000
ExecStart=/usr/bin/chromium --no-sandbox --start-maximized --disable-gpu --disable-dev-shm-usage
Restart=always
User=wilson_austin10
Group=wilson_austin10
WorkingDirectory=/home/wilson_austin10

[Install]
WantedBy=multi-user.target
